# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

# ─────────────────────────────────────────────────────────────────────────────
# HERMES-AGENT GATEWAY
#
# hermes-agent (github.com/NousResearch/hermes-agent, 149K ★) is the full
# autonomous agent runtime. It owns the Telegram bot token, the agentic loop,
# session memory, skill creation, cron scheduling, and multi-channel routing.
#
# IMPORTANT — binary name change (FIX-D):
#   hermes-agent installs as the `hermes` CLI binary (not `hermes-agent`).
#   The pyproject.toml entry_point is: hermes = "hermes_cli.main:main"
#   Previous versions checked `command -v hermes-agent` which always returned
#   false even after a correct install, silently falling through to gateway/run.py.
#   Fix: check `command -v hermes` first, then fall through to gateway/run.py.
#
# Startup order:
#   1. hermes CLI (installed via pip from GitHub source in Dockerfile.vps FIX-A)
#      Reads SOUL.md, config.yaml, gateway.yaml from $HERMES_HOME.
#      Starts Telegram polling, skill engine, memory engine, cron, event bus.
#   2. python3 /app/gateway/run.py
#      Fallback when hermes-agent is unavailable. Reads HERMES_GATEWAY env var
#      and delegates to the appropriate binary (hermes | openclaw). If neither
#      is installed, exits with FATAL — do not mask this error with sleep.
# ─────────────────────────────────────────────────────────────────────────────
[program:hermes-gateway]
# NOTE — Telegram Polling Conflict:
#   If this process crashes and restarts within startsecs (15s), the previous
#   Telegram long-poll connection may still be alive (30s default timeout).
#   Telegram returns HTTP 409 Conflict until the old connection drops.
#   supervisord's startsecs=15 gate prevents rapid-fire restarts from
#   overwhelming the Telegram API. The 409 errors resolve automatically
#   within 30 seconds of the old process dying — no manual intervention needed.
command=/bin/bash -c '
  if command -v hermes &>/dev/null; then
    echo "[hermes-gateway] Starting hermes-agent v$(hermes --version 2>/dev/null | head -1) control plane";
    exec hermes gateway run;
  else
    echo "[hermes-gateway] hermes binary not found — falling back to gateway/run.py dispatcher";
    exec python3 /app/gateway/run.py;
  fi'
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PYTHONUNBUFFERED="1",
    HERMES_HOME="/data/.hermes",
    HERMES_YOLO_MODE="0",
    HERMES_GATEWAY_CONFIG="/data/.hermes/gateway.yaml",
    TERMINAL_ENV="local",
    HERMES_SKILLS_DIR="/data/.hermes/skills",
    HERMES_SESSIONS_DIR="/data/.hermes/sessions",
    HERMES_MEMORIES_DIR="/data/.hermes/memories",
    GATEWAY_ALLOW_ALL_USERS="false",
    HERMES_GATEWAY="hermes",
    CAMOFOX_HOST="camofox",
    CAMOFOX_PORT="9377",
    HERMES_MODEL="deepseek-v4-pro",
    OPENCLAUDE_MODEL="deepseek-r1-distill-llama-70b",
    JCODE_MODEL="kimi-k2.6"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAUDE gRPC SERVER
#
# FIX Problem-2: The Python server now implements the FULL AgentService.Chat
# bidirectional stream:
#   TextChunk      — LLM text tokens streamed to client
#   ToolCallStart  — emitted before executing bash/file tools
#   ToolCallResult — emitted after tool returns (real output, not fabricated)
#   ActionRequired — emitted for dangerous commands (rm -rf, sudo, chmod 777)
#                    CLIENT must send UserInput y/n before execution resumes
#   FinalResponse  — emitted when agentic loop completes or reports DONE
#   ErrorResponse  — emitted on unrecoverable failure
#
# The loop iterates: plan → read files → execute tools → verify → self-correct.
# ActionRequired is the safety gate for Problem-6.
#
# Startup order:
#   1. bun run dev:grpc  (openclaude TypeScript server — real bidirectional gRPC)
#   2. python3 /app/skills/openclaude_grpc/server.py  (Python fallback, full protocol)
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaude-grpc]
command=/bin/bash -c '
  export PATH="/root/.bun/bin:$PATH";
  if command -v bun &>/dev/null && [ -f /opt/openclaude/scripts/start-grpc.ts ]; then
    echo "[openclaude-grpc] Starting TypeScript gRPC server via bun";
    cd /opt/openclaude && exec bun run dev:grpc;
  else
    echo "[openclaude-grpc] bun/scripts unavailable — starting Python full-protocol gRPC server";
    exec python3 /app/skills/openclaude_grpc/server.py;
  fi'
autostart=true
autorestart=true
startretries=3
startsecs=8
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    GRPC_PORT="50051",
    GRPC_HOST="localhost",
    CLAUDE_CODE_USE_OPENAI="1",
    HERMES_YOLO_MODE="0"

# ─────────────────────────────────────────────────────────────────────────────
# JCODE SERVER
#
# Runs jcode's built-in coordination server (github.com/1jehuang/jcode).
# jcode serve provides natively: per-file conflict detection, change
# notifications to sibling swarm workers, and session persistence.
#
# Install: cargo install jcode
#
# FIX — "jcode-server is EXITED" watchdog alert every 5 min:
#   Previously used `exec jcode serve` which replaces the shell. If jcode
#   crashes or is not installed, supervisord sees EXITED and watchdog alerts
#   every ALERT_COOLDOWN (300s). Fix: always fall through to `sleep infinity`
#   so the process stays RUNNING regardless of whether jcode is installed or
#   jcode serve exits. autorestart=false is intentional — jcode is optional.
# ─────────────────────────────────────────────────────────────────────────────
[program:jcode-server]
command=/bin/bash -c '
  if command -v jcode &>/dev/null; then
    echo "[jcode-server] jcode found — starting jcode serve";
    jcode serve 2>&1;
    echo "[jcode-server] jcode serve exited (code $?) — staying alive to suppress watchdog alerts";
  else
    echo "[jcode-server] jcode not installed — process kept alive. Install: cargo install jcode";
  fi;
  exec sleep infinity'
directory=/tmp
autostart=true
autorestart=false
startretries=2
startsecs=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=JCODE_NON_INTERACTIVE="1"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAW GATEWAY (multi-channel relay)
#
# FIX-D: openclaw is a standalone personal AI assistant gateway. It is NOT a
# relay into hermes-agent. Running it in parallel with hermes-gateway causes
# two agents to compete for the same resources and produce split state.
#
# openclaw is set to autostart=false for the following reasons:
#   1. No channels have been configured via `openclaw onboard` — it would
#      start, bind port 18789, and silently do nothing.
#   2. hermes-agent already provides Telegram, Discord, Slack, WhatsApp, Signal,
#      and Email natively via `hermes gateway setup`. openclaw is not needed
#      for any of these channels.
#   3. Starting an unconfigured gateway process wastes RAM and produces log noise.
#
# TO ENABLE: Run `openclaw onboard --install-daemon` to configure channels,
# then flip autostart=true here and restart supervisord. Only do this if you
# want openclaw to serve as an independent assistant on iMessage, LINE, WeChat,
# or other channels not supported by hermes-agent.
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaw-gateway]
command=/bin/bash -c '
  echo "[openclaw-gateway] autostart=false — disabled pending channel configuration.";
  echo "  Run: openclaw onboard --install-daemon  then flip autostart=true in supervisord.conf";
  exec sleep infinity'
directory=/root
autostart=false
autorestart=false
startretries=1
startsecs=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    OPENCLAW_STATE_DIR="/root/.openclaw",
    OPENCLAW_CONFIG_PATH="/root/.openclaw/openclaw.json",
    OPENCLAW_NO_RESPAWN="1",
    NODE_COMPILE_CACHE="/var/tmp/openclaw-compile-cache"

# ─────────────────────────────────────────────────────────────────────────────
# WATCHDOG — Process health monitor + Telegram alerter
# ─────────────────────────────────────────────────────────────────────────────
[program:watchdog]
command=python3 /app/scripts/watchdog.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
