# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

# ─────────────────────────────────────────────────────────────────────────────
# HERMES-AGENT GATEWAY
#
# FIX Problem-1: hermes-agent IS the control plane. It owns session routing,
# memory, multi-provider failover, tool calling, and skills — all natively.
# Running gateway/run.py in parallel was architecture theater (Problem 5).
#
# Startup order:
#   1. hermes-agent CLI (pip entry_point — installed by requirements.txt)
#      hermes-agent --config $HERMES_HOME/config.yaml
#      Reads SOUL.md, config.yaml, gateway.yaml. Starts Telegram polling,
#      skill engine, memory engine, cron, event bus, task queue.
#   2. python3 /app/gateway/run.py
#      Fallback when hermes-agent package is unavailable. Same SOUL.md, same
#      tool set, same 4-tier search cascade — but no native hermes infrastructure.
#
# hermes-agent handles everything gateway/run.py re-implemented:
#   Telegram polling, conversation history, tool calling, memory injection,
#   anti-fabrication (via SOUL.md rules), multi-provider routing (transport ABC),
#   session management, cron scheduling, subagent protocol.
# ─────────────────────────────────────────────────────────────────────────────
[program:hermes-gateway]
command=/bin/bash -c '
  if command -v hermes-agent &>/dev/null; then
    echo "[hermes-gateway] Starting hermes-agent control plane";
    exec hermes-agent --config "$HERMES_HOME/config.yaml";
  else
    echo "[hermes-gateway] hermes-agent not found — falling back to gateway/run.py";
    exec python3 /app/gateway/run.py;
  fi'
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PYTHONUNBUFFERED="1",
    HERMES_HOME="/data/.hermes",
    HERMES_YOLO_MODE="1",
    HERMES_GATEWAY_CONFIG="/data/.hermes/gateway.yaml",
    TERMINAL_ENV="local",
    HERMES_SKILLS_DIR="/data/.hermes/skills",
    HERMES_SESSIONS_DIR="/data/.hermes/sessions",
    HERMES_MEMORIES_DIR="/data/.hermes/memories",
    GATEWAY_ALLOW_ALL_USERS="true",
    CAMOFOX_HOST="camofox",
    CAMOFOX_PORT="9377",
    HERMES_MODEL="deepseek-v4-pro",
    OPENCLAUDE_MODEL="deepseek-r1-distill-llama-70b",
    JCODE_MODEL="kimi-k2.6"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAUDE gRPC SERVER
#
# FIX Problem-2: The Python server now implements the FULL AgentService.Chat
# bidirectional stream:
#   TextChunk      — LLM text tokens streamed to client
#   ToolCallStart  — emitted before executing bash/file tools
#   ToolCallResult — emitted after tool returns (real output, not fabricated)
#   ActionRequired — emitted for dangerous commands (rm -rf, sudo, chmod 777)
#                    CLIENT must send UserInput y/n before execution resumes
#   FinalResponse  — emitted when agentic loop completes or reports DONE
#   ErrorResponse  — emitted on unrecoverable failure
#
# The loop iterates: plan → read files → execute tools → verify → self-correct.
# ActionRequired is the safety gate for Problem-6.
#
# Startup order:
#   1. bun run dev:grpc  (openclaude TypeScript server — real bidirectional gRPC)
#   2. python3 /app/skills/openclaude_grpc/server.py  (Python fallback, full protocol)
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaude-grpc]
command=/bin/bash -c '
  export PATH="/root/.bun/bin:$PATH";
  if command -v bun &>/dev/null && [ -f /opt/openclaude/scripts/start-grpc.ts ]; then
    echo "[openclaude-grpc] Starting TypeScript gRPC server via bun";
    cd /opt/openclaude && exec bun run dev:grpc;
  else
    echo "[openclaude-grpc] bun/scripts unavailable — starting Python full-protocol gRPC server";
    exec python3 /app/skills/openclaude_grpc/server.py;
  fi'
autostart=true
autorestart=true
startretries=3
startsecs=8
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    GRPC_PORT="50051",
    GRPC_HOST="localhost",
    CLAUDE_CODE_USE_OPENAI="1",
    HERMES_YOLO_MODE="0"

# ─────────────────────────────────────────────────────────────────────────────
# JCODE SERVER (coordination daemon)
#
# FIX Problem-3: jcode serve now actually runs a coordination server at :7865.
# The Rhodawk jcode wrapper now checks for a real jcode binary first, and
# falls back to the API wrapper only when the binary is not installed.
# The coordination daemon tracks per-file locks so parallel swarm workers
# cannot write conflicting edits to the same file.
# ─────────────────────────────────────────────────────────────────────────────
[program:jcode-server]
command=/bin/bash -c '
  if command -v jcode &>/dev/null && jcode --version 2>/dev/null | grep -v rhodawk | grep -q "jcode"; then
    echo "[jcode-server] Starting real jcode server on port 7865";
    exec jcode serve --port 7865 2>&1;
  else
    echo "[jcode-server] Real jcode binary not found — starting Rhodawk coordination daemon";
    exec python3 /app/skills/jcode_swarm/coord_daemon.py --port 7865;
  fi'
directory=/tmp
autostart=true
autorestart=false
startretries=2
startsecs=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=JCODE_NON_INTERACTIVE="1"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAW GATEWAY (multi-channel relay)
#
# FIX Problem-4: TELEGRAM_BOT_TOKEN is explicitly unset before openclaw starts.
# This prevents openclaw's built-in Telegram adapter from fighting with
# hermes-gateway over the same bot token (both would poll the same token,
# each getting ~50% of messages, both handling them in isolation — broken).
#
# openclaw runs as a multi-channel RELAY (WhatsApp, Discord, Slack, Signal,
# iMessage, 16+ other channels) on port 18789 — distinct from hermes-gateway.
# It bridges those channels INTO hermes-gateway via the ACP subagent protocol,
# not as a competing Telegram bot.
#
# The bash retry loop keeps supervisord happy (RUNNING) while handling crashes.
# autorestart=false: the loop manages its own restarts.
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaw-gateway]
command=/bin/bash -c '
  unset TELEGRAM_BOT_TOKEN;
  if ! command -v openclaw &>/dev/null; then
    echo "[openclaw-gateway] not found — sleeping (non-critical)";
    exec sleep infinity;
  fi;
  while true; do
    echo "[openclaw-gateway] launching openclaw gateway --port 18789";
    openclaw gateway --port 18789 2>&1;
    echo "[openclaw-gateway] exited (rc=$?) — retrying in 30s";
    sleep 30;
  done'
directory=/root
autostart=true
autorestart=false
startretries=1
startsecs=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    OPENCLAW_STATE_DIR="/root/.openclaw",
    OPENCLAW_CONFIG_PATH="/root/.openclaw/openclaw.json",
    OPENCLAW_NO_RESPAWN="1",
    NODE_COMPILE_CACHE="/var/tmp/openclaw-compile-cache"

# ─────────────────────────────────────────────────────────────────────────────
# WATCHDOG — Process health monitor + Telegram alerter
# ─────────────────────────────────────────────────────────────────────────────
[program:watchdog]
command=python3 /app/scripts/watchdog.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
