# supervisord.peak.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

# ─────────────────────────────────────────────────────────────────────────────
# HERMES-GATEWAY
# Role:    Full orchestration engine + Telegram adapter (NousResearch)
#          Peak architecture entry point: handles all inbound messages,
#          routes to skills, manages sessions, injects memory context.
#
# Command: python3 /app/gateway/run.py (production gateway)
# ─────────────────────────────────────────────────────────────────────────────
[program:hermes-gateway]
command=python3 /app/gateway/run.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PYTHONUNBUFFERED="1",
    HERMES_HOME="/data/.hermes",
    HERMES_YOLO_MODE="1",
    HERMES_GATEWAY_CONFIG="/data/.hermes/gateway.yaml",
    TERMINAL_ENV="local",
    HERMES_SKILLS_DIR="/data/.hermes/skills",
    HERMES_SESSIONS_DIR="/data/.hermes/sessions",
    HERMES_MEMORIES_DIR="/data/.hermes/memories",
    GATEWAY_ALLOW_ALL_USERS="true",
    CAMOFOX_HOST="camofox",
    CAMOFOX_PORT="9377",
    REDIS_URL="redis://redis:6379/0",
    HERMES_MODEL="deepseek-v4-pro",
    OPENCLAUDE_MODEL="deepseek-r1-distill-llama-70b",
    JCODE_MODEL="kimi-k2.6"

# ─────────────────────────────────────────────────────────────────────────────
# EVENT ROUTER
# Role:    Consumes events from Redis PubSub and dispatches to handlers.
#          Part of the peak architecture event-driven backbone.
#          Bridges external webhooks and internal actions to the gateway.
# ─────────────────────────────────────────────────────────────────────────────
[program:event-router]
command=python3 /app/scripts/event_worker.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PYTHONUNBUFFERED="1",
    REDIS_URL="redis://redis:6379/1",
    EVENT_BUS_CHANNELS="hermes:events,hermes:tasks,hermes:alerts,hermes:health"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAUDE gRPC SERVER
# Role:    Persistent gRPC coding server at localhost:50051
#          Enables bidirectional streaming, action_required callbacks,
#          and programmatic sub-agent control.
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaude-grpc]
command=/bin/bash -c 'export PATH="/root/.bun/bin:$PATH"; if command -v bun &>/dev/null && [ -f /opt/openclaude/scripts/start-grpc.ts ]; then cd /opt/openclaude && exec bun run dev:grpc; else echo "[openclaude-grpc] bun dev:grpc unavailable — using Python gRPC server"; exec python3 /app/openclaude_grpc/server.py; fi'
autostart=true
autorestart=true
startretries=3
startsecs=8
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    GRPC_PORT="50051",
    GRPC_HOST="localhost",
    CLAUDE_CODE_USE_OPENAI="1"

# ─────────────────────────────────────────────────────────────────────────────
# TASK WORKERS
# Role:    Background task execution pool (dramatiq-style workers).
#          Dequeues tasks from Redis sorted set, executes via skill handlers,
#          reports progress/completion/failure back to the event bus.
# ─────────────────────────────────────────────────────────────────────────────
[program:task-workers]
command=python3 /app/scripts/task_worker.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=25
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PYTHONUNBUFFERED="1",
    REDIS_URL="redis://redis:6379/0",
    TASK_WORKERS="5",
    TASK_DEFAULT_TIMEOUT="3600"

# ─────────────────────────────────────────────────────────────────────────────
# JCODE SERVER (background)
# Role:    Serves the jcode agent runtime on localhost:7865
#          hermes-agent calls `jcode run --message "..."` for scaffolding tasks.
#          Non-critical — gateway continues if jcode fails to start.
# ─────────────────────────────────────────────────────────────────────────────
[program:jcode-server]
command=/bin/bash -c 'if command -v jcode &>/dev/null; then jcode serve --port 7865 2>&1; else echo "[jcode-server] not installed — sleeping (non-critical)"; exec sleep infinity; fi'
directory=/tmp
autostart=true
autorestart=false
startretries=2
startsecs=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=JCODE_NON_INTERACTIVE="1"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAW GATEWAY (multi-channel relay)
# Role:    Personal AI assistant gateway at localhost:18789
#          Adds WhatsApp, Discord, Slack, Signal, and other channels.
#          Non-critical — hermes-gateway and openclaude-grpc continue
#          if openclaw fails to start.
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaw-gateway]
command=/bin/bash -c 'unset TELEGRAM_BOT_TOKEN; if ! command -v openclaw &>/dev/null; then echo "[openclaw-gateway] not found -- sleeping"; exec sleep infinity; fi; while true; do echo "[openclaw-gateway] launching openclaw gateway --port 18789"; openclaw gateway --port 18789 2>&1; echo "[openclaw-gateway] exited -- retrying in 30s"; sleep 30; done'
directory=/root
autostart=true
autorestart=false
startretries=1
startsecs=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    OPENCLAW_STATE_DIR="/root/.openclaw",
    OPENCLAW_CONFIG_PATH="/root/.openclaw/openclaw.json",
    OPENCLAW_NO_RESPAWN="1",
    NODE_COMPILE_CACHE="/var/tmp/openclaw-compile-cache"

# ─────────────────────────────────────────────────────────────────────────────
# WATCHDOG — Process health monitor + Telegram alerter
# Role:    Polls supervisord every 60s. Sends Telegram messages to the owner
#          if any process enters FATAL or EXITED. Sends recovery notices when
#          services come back to RUNNING.
# ─────────────────────────────────────────────────────────────────────────────
[program:watchdog]
command=python3 /app/scripts/watchdog.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
