# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

# ─────────────────────────────────────────────────────────────────────────────
# HERMES-AGENT GATEWAY
# Role:    Full orchestration engine + Telegram adapter (NousResearch)
#          Uses hermes-agent's built-in gateway/platforms/telegram.py.
#          AIAgent handles: LLM routing, tool calling, skills, memory,
#          context compression, provider failover, session management.
#
# Command resolution order:
#   1. hermes-agent CLI (pip entry_point — preferred)
#   2. python3 -m hermes_agent (module fallback — same package)
#   3. python3 /app/main.py (final fallback — supervised restart loop)
#
# FIX: was `python3 -m gateway.run` which is NOT the correct module path.
#      hermes-agent installs as the `hermes_agent` Python package.
#      The CLI entry point `hermes-agent` is the canonical way to start it.
# ─────────────────────────────────────────────────────────────────────────────
[program:hermes-gateway]
command=python3 /app/gateway/run.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PYTHONUNBUFFERED="1",
    HERMES_HOME="/data/.hermes",
    HERMES_YOLO_MODE="1",
    HERMES_GATEWAY_CONFIG="/data/.hermes/gateway.yaml",
    TERMINAL_ENV="local",
    HERMES_SKILLS_DIR="/data/.hermes/skills",
    HERMES_SESSIONS_DIR="/data/.hermes/sessions",
    HERMES_MEMORIES_DIR="/data/.hermes/memories",
    GATEWAY_ALLOW_ALL_USERS="true",
    CAMOFOX_HOST="camofox",
    CAMOFOX_PORT="9377",
    HERMES_MODEL="deepseek-v4-pro",
    OPENCLAUDE_MODEL="deepseek-r1-distill-llama-70b",
    JCODE_MODEL="kimi-k2.6"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAUDE gRPC SERVER
# Role:    Persistent gRPC coding server at localhost:50051
#          Enables bidirectional streaming, action_required callbacks,
#          and programmatic sub-agent control.
#
# Command resolution order:
#   1. npm run dev:grpc  (from @gitlawb/openclaude package dir)
#   2. openclaude --grpc (CLI flag fallback)
#   3. sleep infinity    (graceful no-op if binary unavailable)
#
# FIX: added `sleep infinity` fallback so supervisor does not crash-loop if
#      openclaude is unavailable, keeping hermes-gateway alive.
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaude-grpc]
command=/bin/bash -c 'export PATH="/root/.bun/bin:$PATH"; if command -v bun &>/dev/null && [ -f /opt/openclaude/scripts/start-grpc.ts ]; then cd /opt/openclaude && exec bun run dev:grpc; else echo "[openclaude-grpc] bun dev:grpc unavailable — using Python gRPC server"; exec python3 /app/skills/openclaude_grpc/server.py; fi'
autostart=true
autorestart=true
startretries=3
startsecs=8
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    GRPC_PORT="50051",
    GRPC_HOST="localhost",
    CLAUDE_CODE_USE_OPENAI="1"

# ─────────────────────────────────────────────────────────────────────────────
# JCODE SERVER (background)
# Role:    Serves the jcode agent runtime on localhost:7865
#          hermes-agent calls `jcode run --message "..."` for scaffolding tasks
#          The server stays warm between calls for lower latency.
#          Non-critical — gateway continues if jcode fails to start.
# ─────────────────────────────────────────────────────────────────────────────
[program:jcode-server]
command=/bin/bash -c 'if command -v jcode &>/dev/null; then jcode serve --port 7865 2>&1; else echo "[jcode-server] not installed — sleeping (non-critical)"; exec sleep infinity; fi'
directory=/tmp
autostart=true
autorestart=false
startretries=2
startsecs=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=JCODE_NON_INTERACTIVE="1"

# ─────────────────────────────────────────────────────────────────────────────
# OPENCLAW GATEWAY (multi-channel relay)
# Role:    Personal AI assistant gateway at localhost:18789
#          Adds WhatsApp, Discord, Slack, Signal, iMessage, and 20+ channels
#          without touching hermes-gateway. Non-critical — hermes-gateway and
#          openclaude-grpc continue if openclaw fails to start.
# Provider: DO Inference via models.providers.do-inference (written by init)
#
# DESIGN: The command is a bash retry loop, NOT a direct exec.
#   The bash wrapper process never exits — supervisord always sees RUNNING.
#   If openclaw crashes (any reason: bad flag, missing config, port conflict),
#   the loop logs the exit code, waits 30s, and retries automatically.
#   This prevents FATAL regardless of why openclaw exits.
#   startsecs=3 matches bash startup time (not openclaw startup time).
#   autorestart=false: the loop handles its own restarts — supervisord
#   does not need to restart the bash wrapper (it never exits).
# ─────────────────────────────────────────────────────────────────────────────
[program:openclaw-gateway]
command=/bin/bash -c 'unset TELEGRAM_BOT_TOKEN; if ! command -v openclaw &>/dev/null; then echo "[openclaw-gateway] not found -- sleeping"; exec sleep infinity; fi; while true; do echo "[openclaw-gateway] launching openclaw gateway --port 18789"; openclaw gateway --port 18789 2>&1; echo "[openclaw-gateway] exited -- retrying in 30s"; sleep 30; done'
directory=/root
autostart=true
autorestart=false
startretries=1
startsecs=3
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    OPENCLAW_STATE_DIR="/root/.openclaw",
    OPENCLAW_CONFIG_PATH="/root/.openclaw/openclaw.json",
    OPENCLAW_NO_RESPAWN="1",
    NODE_COMPILE_CACHE="/var/tmp/openclaw-compile-cache"

# ─────────────────────────────────────────────────────────────────────────────
# WATCHDOG — Process health monitor + Telegram alerter
# Role:    Polls supervisord every 60s. Sends a Telegram message to the owner
#          if any process enters FATAL or EXITED. Sends a recovery notice when
#          it comes back to RUNNING. Uses only Python stdlib — no extra deps.
# ─────────────────────────────────────────────────────────────────────────────
[program:watchdog]
command=python3 /app/scripts/watchdog.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
